# secretary-bot.service.example
# secretary-bot.service.example
#
# Пример systemd unit-файла для secretary-bot.
# Скопируйте в /etc/systemd/system/secretary-bot.service,
# отредактируйте пути при необходимости, затем выполните:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now secretary-bot

[Unit]
Description=Secretary Bot — Telegram Business assistant
After=network.target

[Service]
Type=simple
User=mkbot
Group=mkbot
WorkingDirectory=/home/mkbot/secretary-bot
ExecStart=/home/mkbot/secretary-bot/start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
